# misc/sandbox-image/neko/chromium.conf
[program:chromium]
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",DISPLAY="%(ENV_DISPLAY)s"
; The profile lives on the PVC and is reused across sandboxes for the same
; user; a previous (terminated) sandbox can leave a Chromium SingletonLock that
; makes a new Chromium refuse to start ("profile in use"). Remove the stale lock
; before each (re)start — only one sandbox uses the profile at a time.
command=/bin/sh -c "rm -f /workspace/.cubebox-browser-profile/Singleton* 2>/dev/null; exec /ms-playwright/chrome --no-sandbox --test-type --window-position=0,0 --display=%(ENV_DISPLAY)s --user-data-dir=/workspace/.cubebox-browser-profile --no-first-run --no-default-browser-check --disable-infobars --start-maximized --force-dark-mode --disable-gpu --disable-software-rasterizer --disable-dev-shm-usage --remote-debugging-port=9222 --remote-debugging-address=127.0.0.1 --remote-allow-origins=* about:blank"
stopsignal=INT
autorestart=true
priority=800
user=%(ENV_USER)s
stdout_logfile=/var/log/neko/chromium.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
redirect_stderr=true

[program:openbox]
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",DISPLAY="%(ENV_DISPLAY)s"
command=/usr/bin/openbox --config-file /etc/neko/openbox.xml
autorestart=true
priority=300
user=%(ENV_USER)s
stdout_logfile=/var/log/neko/openbox.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
redirect_stderr=true
